# data/rpi-pipackage.service.in
[Unit]
Description=Raspberry Pi Package Installer Service

[Service]
ExecStart=@servicedir@/pi-gpk-dbus-service
Type=dbus
BusName=org.rpi.pipackage

[Install]
WantedBy=default.target
